# vpn-bot.service
[Unit]
Description=VPN Manager Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vpn-installer
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/vpn-installer/vpn_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 